# frp/frp.conf
[program:frpc]
command=/home/diggerdu/bin/frpc -c frpc.ini
directory=/home/diggerdu/bin/
autorestart=true
autostart=true

