# frp/frp.conf
[program:frpc]
command=./frpc -c frpc.ini
autorestart=true
autostart=true
[supervisord]
